# config/supervisord_celery.conf
[program:celery]
directory=/rapidpro
command=python manage.py celery worker -B --loglevel=WARNING --autoscale=3,1
#startsecs=0
redirect_stderr=true
autorestart=true
